# tools/launchd/com.ecoplating.steelhead-escalation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ecoplating.steelhead-escalation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/oviazcan/Projects/Ecoplating/SteelheadAutomator/tools/run-escalation.sh</string>
    </array>
    <!-- A :53, 30 min después del motor (:23), para que el needs-attention del ciclo
         ya esté escrito. Gate por needs-attention.json → casi siempre sale en <1s. -->
    <key>StartCalendarInterval</key>
    <dict><key>Minute</key><integer>53</integer></dict>
    <key>WorkingDirectory</key>
    <string>/Users/oviazcan/Projects/Ecoplating/SteelheadAutomator</string>
    <key>StandardOutPath</key>
    <string>/Users/oviazcan/Projects/Ecoplating/SteelheadAutomator/tools/.hash-autopilot/escalation.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/oviazcan/Projects/Ecoplating/SteelheadAutomator/tools/.hash-autopilot/escalation.err.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
